# itamae/cookbooks/common/acmesmith-fetch/files/etc/systemd/system/nkmi-acmesmith-fetch.service
[Unit]
Description=acmesmith-fetch

[Service]
Type=oneshot
ExecStart=/usr/bin/nkmi-acmesmith-fetch
