# itamae/cookbooks/arch/etcd/templates/etc/systemd/system/etcd.service.d/env.conf
[Service]
Environment=
<%- node[:etcd][:server][:env].each do |k,v| -%>
Environment=ETCD_<%= k.upcase %>=<%= v %>
<%- end -%>
