# itamae/cookbooks/common/vault-renew/files/etc/systemd/system/nkmi-vault-renew.service
[Unit]
Description=Renew periodic tokens

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/nkmi-vault-renew
