# itamae/cookbooks/common/vault-cert/files/etc/systemd/system/nkmi-vault-cert.service
[Unit]
Description=Retrieve certificates from Vault PKI

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/nkmi-vault-cert
