# itamae/cookbooks/common/vault-cert/files/etc/systemd/system/nkmi-vault-cert.timer
[Unit]
Description=periodically renew vault pki certificates

[Timer]
OnBootSec=15min
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target
