# itamae/cookbooks/arch/docker/templates/etc/systemd/system/docker.service.d/dependencies.conf
# vim: ft=systemd
[Unit]
Wants=nkmi-mnt-vol.service
After=nkmi-mnt-vol.service
