# itamae/cookbooks/common/ec2-vxlan/files/etc/systemd/system/nkmi-ec2-vxlan-auto.service
[Unit]
Description=Ensure VXLAN interface and fdb up to date

[Service]
Type=oneshot
ExecStart=/usr/bin/nkmi-ec2-vxlan-auto
